# deploy/gunicorn.service
[Unit]
Description=Gunicorn daemon for Stael Di
After=network.target

[Service]
User=door
Group=www-data
WorkingDirectory=/home/door/site
ExecStart=/home/door/site/venv/bin/gunicorn door.wsgi:application \
    --bind 127.0.0.1:8000 \
    --workers 3 \
    --timeout 120 \
    --access-logfile /var/log/gunicorn/access.log \
    --error-logfile /var/log/gunicorn/error.log
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
